# docs/legacy/transcribe-daemon.service
[Unit]
Description=Transcribe Daemon (Parakeet ASR)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/transcribe-daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
