# pkg/parakeet-transcribe/usr/lib/systemd/user/parakeet-transcribe.service
[Unit]
Description=Parakeet Transcribe Daemon (ASR)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/transcribe-daemon /usr/share/parakeet-transcribe/tdt
Restart=on-failure
RestartSec=5
Environment=XDG_RUNTIME_DIR=/run/user/%U

[Install]
WantedBy=default.target
