# pkg/dictee/DEBIAN/postrm
#!/bin/bash
set -e

if [ "$1" = "remove" ] || [ "$1" = "purge" ]; then
    # Arrêter et désactiver le service pour tous les utilisateurs connectés
    for uid in $(loginctl list-sessions --no-legend 2>/dev/null | awk '{print $2}' | sort -u); do
        user=$(id -nu "$uid" 2>/dev/null) || continue
        if [ -d "/run/user/$uid" ]; then
            sudo -u "$user" XDG_RUNTIME_DIR="/run/user/$uid" DBUS_SESSION_BUS_ADDRESS="unix:path=/run/user/$uid/bus" \
                systemctl --user disable --now dictee-tray 2>/dev/null || true
            sudo -u "$user" XDG_RUNTIME_DIR="/run/user/$uid" DBUS_SESSION_BUS_ADDRESS="unix:path=/run/user/$uid/bus" \
                systemctl --user disable --now dictee 2>/dev/null || true
        fi
    done
fi

exit 0
