# pkg/dictee/DEBIAN/postrm
#!/bin/bash
set -e

if [ "$1" = "remove" ] || [ "$1" = "purge" ]; then
    # Stop and disable all user services
    for uid in $(loginctl list-sessions --no-legend 2>/dev/null | awk '{print $2}' | sort -u); do
        user=$(id -nu "$uid" 2>/dev/null) || continue
        [ "$user" = "root" ] && continue
        [ -d "/run/user/$uid" ] || continue
        _run="sudo -u $user XDG_RUNTIME_DIR=/run/user/$uid DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/$uid/bus"
        $_run systemctl --user stop dictee-ptt dictee-tray dotoold dictee dictee-vosk dictee-whisper dictee-canary 2>/dev/null || true
        $_run systemctl --user disable dictee-ptt dictee-tray dotoold dictee dictee-vosk dictee-whisper dictee-canary 2>/dev/null || true
        $_run systemctl --user daemon-reload 2>/dev/null || true
    done

    # Clean locales
    for lang in fr de es it uk pt; do
        rm -f "/usr/share/locale/$lang/LC_MESSAGES/dictee.mo"
    done

    # Reload udev rules
    udevadm control --reload-rules 2>/dev/null || true
    udevadm trigger /dev/uinput 2>/dev/null || true

    # Remove venv Python (text2num) — créé au postinst, pas tracké par dpkg
    rm -rf /usr/share/dictee/postprocess-env

    # Remove venv CUDA (nvidia-*-cu12 wheels) — créé au postinst de dictee-cuda
    rm -rf /opt/dictee/cuda-venv
    # Clean up dangling symlinks in /usr/lib/dictee/ (pointing to the removed venv)
    if [ -d /usr/lib/dictee ]; then
        find /usr/lib/dictee -maxdepth 1 -type l ! -exec test -e '{}' \; -delete 2>/dev/null || true
    fi
    # Remove empty /opt/dictee if we created it
    rmdir /opt/dictee 2>/dev/null || true
fi

exit 0
