# pkg/dictee/usr/lib/systemd/user/dictee-tray.service
[Unit]
Description=Parakeet Transcribe Tray Icon
After=dictee.service
BindsTo=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/dictee-tray
Restart=on-failure
RestartSec=3

[Install]
WantedBy=graphical-session.target
