# pkg/dictee/usr/lib/systemd/user/dictee.service
[Unit]
Description=Parakeet Transcribe Daemon (ASR)
After=network.target
Conflicts=dictee-vosk.service dictee-whisper.service dictee-canary.service

[Service]
Type=simple
ExecStart=/usr/bin/transcribe-daemon
ExecStartPost=/bin/bash -c 'echo idle > /dev/shm/.dictee_state'
ExecStopPost=/bin/bash -c 'echo offline > /dev/shm/.dictee_state'
Restart=on-failure
RestartSec=5
EnvironmentFile=-%h/.config/dictee.conf
Environment=XDG_RUNTIME_DIR=/run/user/%U
Environment=ORT_DYLIB_PATH=/usr/lib/dictee/libonnxruntime.so

[Install]
WantedBy=default.target
